# data/ayatana-indicator-messages.service.in
[Unit]
Description=Ayatana Indicator Messages Service
PartOf=graphical-session.target
PartOf=ayatana-indicators.target

[Service]
ExecStart=@pkglibexecdir@/ayatana-indicator-messages-service
Restart=on-failure

[Install]
WantedBy=ayatana-indicators.target
